# codes/Orin/dpkg_build/ota/ota-control-server-1.0.5_arm64/DEBIAN/postinst
#!/bin/bash



systemctl daemon-reload
systemctl enable ota-control.service

systemctl start  ota-control.service
systemctl restart  ota-control.service

exit 0
