# codes/Orin/http_server/http_server/http_upload.service
[Unit]
Description=HTTP Server for Screenshot Upload
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/server.py
WorkingDirectory=/path/to
Restart=always
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target

