# codes/Orin/dpkg_build/upload/upload-file-server-1.0.2_arm64/DEBIAN/postinst
#!/bin/bash



systemctl daemon-reload
systemctl enable upload-file.service

systemctl start  upload-file.service
systemctl restart  upload-file.service

exit 0
